# var/lib/dpkg/info/dphys-swapfile.conffiles
/etc/dphys-swapfile
/etc/bash_completion.d/dphys-swapfile
/etc/init.d/dphys-swapfile
