# var/lib/dpkg/info/isc-dhcp-server.conffiles
/etc/dhcp/dhcpd.conf
/etc/init.d/isc-dhcp-server
