# var/lib/dpkg/info/pigpio.list
/.
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/pigpiod.service
/opt
/opt/pigpio
/opt/pigpio/cgi
/usr
/usr/lib
/usr/lib/libpigpiod_if.so
/usr/lib/libpigpiod_if2.so
/usr/lib/libpigpio.so
/usr/share
/usr/share/doc
/usr/share/doc/pigpio
/usr/share/doc/pigpio/copyright
/usr/share/doc/pigpio/README.gz
/usr/share/doc/pigpio/changelog.Debian.gz
/usr/include
/usr/include/pigpio.h
/usr/include/pigpiod_if2.h
/usr/include/pigpiod_if.h
/usr/man
/usr/man/man1
/usr/man/man1/pigpiod.1.gz
/usr/man/man1/pigs.1.gz
/usr/man/man1/pig2vcd.1.gz
/usr/man/man3
/usr/man/man3/pigpio.3.gz
/usr/man/man3/pigpiod_if2.3.gz
/usr/man/man3/pigpiod_if.3.gz
/usr/bin
/usr/bin/pigs
/usr/bin/pig2vcd
/usr/bin/pigpiod
