# var/lib/dpkg/info/unattended-upgrades.list
/.
/var
/var/log
/var/log/unattended-upgrades
/usr
/usr/share
/usr/share/locale
/usr/share/locale/sk
/usr/share/locale/sk/LC_MESSAGES
/usr/share/locale/sk/LC_MESSAGES/unattended-upgrades.mo
/usr/share/locale/es
/usr/share/locale/es/LC_MESSAGES
/usr/share/locale/es/LC_MESSAGES/unattended-upgrades.mo
/usr/share/locale/fr
/usr/share/locale/fr/LC_MESSAGES
/usr/share/locale/fr/LC_MESSAGES/unattended-upgrades.mo
/usr/share/locale/da
/usr/share/locale/da/LC_MESSAGES
/usr/share/locale/da/LC_MESSAGES/unattended-upgrades.mo
/usr/share/locale/ja
/usr/share/locale/ja/LC_MESSAGES
/usr/share/locale/ja/LC_MESSAGES/unattended-upgrades.mo
/usr/share/locale/pt
/usr/share/locale/pt/LC_MESSAGES
/usr/share/locale/pt/LC_MESSAGES/unattended-upgrades.mo
/usr/share/locale/ru
/usr/share/locale/ru/LC_MESSAGES
/usr/share/locale/ru/LC_MESSAGES/unattended-upgrades.mo
/usr/share/locale/de
/usr/share/locale/de/LC_MESSAGES
/usr/share/locale/de/LC_MESSAGES/unattended-upgrades.mo
/usr/share/locale/eu
/usr/share/locale/eu/LC_MESSAGES
/usr/share/locale/eu/LC_MESSAGES/unattended-upgrades.mo
/usr/share/locale/gl
/usr/share/locale/gl/LC_MESSAGES
/usr/share/locale/gl/LC_MESSAGES/unattended-upgrades.mo
/usr/share/locale/pt_BR
/usr/share/locale/pt_BR/LC_MESSAGES
/usr/share/locale/pt_BR/LC_MESSAGES/unattended-upgrades.mo
/usr/share/doc
/usr/share/doc/unattended-upgrades
/usr/share/doc/unattended-upgrades/README.md
/usr/share/doc/unattended-upgrades/copyright
/usr/share/doc/unattended-upgrades/changelog.gz
/usr/share/unattended-upgrades
/usr/share/unattended-upgrades/20auto-upgrades
/usr/share/unattended-upgrades/20auto-upgrades-disabled
/usr/share/unattended-upgrades/unattended-upgrade-shutdown
/usr/share/man
/usr/share/man/man8
/usr/share/man/man8/unattended-upgrade.8.gz
/usr/lib
/usr/lib/python3
/usr/lib/python3/dist-packages
/usr/lib/python3/dist-packages/unattended_upgrades-0.1.egg-info
/usr/lib/python3/dist-packages/unattended_upgrades-0.1.egg-info/SOURCES.txt
/usr/lib/python3/dist-packages/unattended_upgrades-0.1.egg-info/dependency_links.txt
/usr/lib/python3/dist-packages/unattended_upgrades-0.1.egg-info/top_level.txt
/usr/lib/python3/dist-packages/unattended_upgrades-0.1.egg-info/PKG-INFO
/usr/bin
/usr/bin/unattended-upgrade
/etc
/etc/logrotate.d
/etc/logrotate.d/unattended-upgrades
/etc/pm
/etc/pm/sleep.d
/etc/pm/sleep.d/10_unattended-upgrades-hibernate
/etc/apt
/etc/apt/apt.conf.d
/etc/apt/apt.conf.d/50unattended-upgrades
/etc/init.d
/etc/init.d/unattended-upgrades
/usr/bin/unattended-upgrades
